# packages/systemd/.config/systemd/user/vllm-qwen.service
[Unit]
Description=vLLM inference server - QuantTrio/Qwen3.6-27B-AWQ (port 8082)
Documentation=https://github.com/vllm-project/vllm
After=network-online.target docker.service
Wants=network-online.target
Conflicts=vllm-gemma4.service vllm-step3.service

[Service]
Type=simple
ExecStartPre=-/usr/bin/docker rm -f vllm-qwen
ExecStart=/usr/bin/docker run \
  --name vllm-qwen \
  --rm \
  --device=nvidia.com/gpu=all \
  -v %h/.cache/huggingface:/root/.cache/huggingface \
  -p 8080:8080 \
  --ipc=host \
  --shm-size=24g \
  vllm/vllm-openai:v0.20.2 \
  QuantTrio/Qwen3.6-27B-AWQ \
  --port 8080 \
  --max-model-len 4096 \
  --gpu-memory-utilization 0.98 \
  --max-num-seqs 4 \
  --reasoning-parser qwen3
ExecStop=/usr/bin/docker stop vllm-qwen
ExecStopPost=/usr/bin/docker rm -f vllm-qwen
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
