# packages/systemd/.config/systemd/user/vllm-step3.service
[Unit]
Description=vLLM inference server - stepfun-ai/Step3-VL-10B (port 8081)
Documentation=https://github.com/vllm-project/vllm
After=network-online.target docker.service
Wants=network-online.target
Conflicts=vllm-gemma4.service vllm-qwen.service

[Service]
Type=simple
ExecStartPre=-/usr/bin/docker rm -f vllm-step3
ExecStart=/usr/bin/docker run \
  --name vllm-step3 \
  --rm \
  --device=nvidia.com/gpu=all \
  -v %h/.cache/huggingface:/root/.cache/huggingface \
  -p 8080:8080 \
  --ipc=host \
  --shm-size=24g \
  vllm/vllm-openai:v0.20.2 \
  stepfun-ai/Step3-VL-10B \
  --port 8080 \
  --trust-remote-code \
  --max-model-len 4096 \
  --gpu-memory-utilization 0.95 \
  --max-num-seqs 4 \
  --reasoning-parser deepseek_r1
ExecStop=/usr/bin/docker stop vllm-step3
ExecStopPost=/usr/bin/docker rm -f vllm-step3
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
